# systemd/ingress.container
[Unit]
Wants=image-bootstrap@%N.service image-build@%N.timer
After=image-bootstrap@%N.service

[Container]
Image=localhost/ingress
Network=host
AutoUpdate=local
Volume=/var/mnt/persist/ingress:/data
Secret=ollama_api_key,type=env,target=OLLAMA_API_KEY
HealthCmd=curl --fail --silent --show-error localhost || exit 1
HealthOnFailure=kill
HealthInterval=5m

[Service]
Restart=on-failure
Slice=prod-tier0.slice

[Install]
WantedBy=prod.target
