# systemd/ingress.container
[Unit]
Wants=image-build@%N.timer

[Container]
Image=localhost/ingress
Network=host
AutoUpdate=local
Volume=ingress:/data

[Service]
Restart=on-failure

[Install]
WantedBy=prod.target
