# systemd/ingress.container
[Unit]
Wants=image-bootstrap@%N.service image-build@%N.timer
After=image-bootstrap@%N.service

[Container]
Image=localhost/ingress
Network=host
AutoUpdate=local
Volume=/var/mnt/persist/ingress:/data
Secret=ollama_api_key,type=env,target=OLLAMA_API_KEY

[Service]
Restart=on-failure

[Install]
WantedBy=prod.target
