# systemd/mail2task.container
[Unit]
After=ollama.service
Wants=ollama.service

[Container]
AutoUpdate=registry
Image=ghcr.io/ykgmfq/mail2task:latest
ContainerName=mail2task
Network=ollama.network
ReadOnly=true
Environment=OLLAMA_HOST=http://ollama:11434
Environment=OLLAMA_MODEL=gemma4:e4b
Environment=IMAP_HOST=imap.strato.de
Environment=IMAP_USER=verwaltung@poepperl.online
Environment=IMAP_FOLDER=Todo
Environment=IMAP_ARCHIVE_FOLDER=Archive
Environment=TODOIST_PROJECT_NAME=Haushalt
Secret=strato_verwaltung_plain,type=env,target=IMAP_PASSWORD
Secret=todoist_api_plain,type=env,target=TODOIST_API_TOKEN

[Service]
Restart=on-failure
Slice=prod-tier1.slice

[Install]
WantedBy=prod.target
