# images/cloud/src/units/cloud-maintenance.timer
[Timer]
OnCalendar=Monthly
RandomizedDelaySec=10min

[Install]
WantedBy=default.target
