# images/cloud/src/units/cloud-previews.timer
[Unit]

[Timer]
OnCalendar=Sun

[Install]
WantedBy=default.target
