# images/cloud/src/units/cloud-maintenance.timer
[Timer]
OnCalendar=Sat *-*-* 01:00:00

[Install]
WantedBy=default.target
