# images/cloud/src/units/cloud-cron.timer
[Unit]
BindsTo=php-fpm.service

[Timer]
OnCalendar=*:0/5

[Install]
WantedBy=default.target
