# images/cloud/src/units/cloud-db-indices.timer
[Timer]
OnCalendar=Sat *-*-1..7
RandomizedDelaySec=10min

[Install]
WantedBy=default.target
