# images/cloud/src/units/cloud-init.service
[Unit]
Before=php-fpm.service

[Service]
User=cloud
Type=oneshot
WorkingDirectory=/usr/share/nextcloud
Slice=cloud-maint.slice
ExecStart=/usr/bin/bash /usr/local/bin/cloud-init.sh
