# images/cloud/src/units/cloud-maintenance.timer
[Timer]
OnCalendar=Sa

[Install]
WantedBy=default.target
