# images/cloud/src/units/container.target
[Unit]
Requires=php-fpm.service
After=php-fpm.service
