# images/freshrss/src/units/actualize.timer
[Unit]
BindsTo=php-fpm.service

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=default.target
