# images/freshrss/src/units/container.target
[Unit]
Requires=php-fpm.service caddy.service
After=php-fpm.service caddy.service
